# supervisord-redbook.conf
[program:redbook]
command=/bin/bash -c "while true; do /usr/bin/php /app/scripts/worker.php ONE_TIME_MODE > /app/log/worker.log; sleep 30; done"
